# main_working_modules/90-hdmi.rules
ACTION=="change", SUBSYSTEM=="HDMI", KERNEL=="hdmi_device", RUN+="/usr/bin/sudo /home/pi/Device_Drivers/main_working_modules/hdmi_udev_script.sh"


